# hw06/monitor-timer.target
[Unit]
Description=SSH monitor target
StopWhenUnneeded=yes
