# hw06/monitor-timer.timer

[Unit]
Description=SSH monitor timer

[Timer]
OnUnitActiveSec=30
OnCalendar=*:*:30
Unit=monitor-timer.target

[Install]
WantedBy=basic.target